# config/rc.local
#!/bin/sh
chdir /home/debian/openag-device-software
/home/debian/openag-device-software/run.sh &
sleep 15
python /home/debian/openag-device-software/scripts/forward_ports.py &
